# awww/awww.service
[Unit]
Description=awww wallpaper daemon
Documentation=https://codeberg.org/lgfae/awww
PartOf=graphical-session.target
After=graphical-session.target
Requisite=graphical-session.target

[Service]
Type=simple
# On start, awww-daemon restores the last wallpaper from its cache (no flash at boot).
ExecStart=/usr/bin/awww-daemon
Restart=on-failure
RestartSec=1

[Install]
WantedBy=graphical-session.target
